# infra/systemd/ml-release-service.live.service
[Unit]
Description=ML Release Manager - Deployed Service (per-release venv)
After=network.target

[Service]
User=serviceuser
WorkingDirectory=/opt/release_manager/current

ExecStartPre=/usr/bin/test -x /opt/release_manager/current/.venv/bin/python
ExecStart=/opt/release_manager/current/.venv/bin/python -m uvicorn service.app:app --host 0.0.0.0 --port 8000

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
