# infra/systemd/ml-release-service.live.service
[Unit]
Description=ML Release Manager - Deployed Service
After=network.target

[Service]
User=serviceuser
WorkingDirectory=/opt/release_manager/current

# Generic FastAPI service entrypoint
ExecStart=/opt/release_manager/venv/bin/uvicorn service.app:app --host 0.0.0.0 --port 8000

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
